# php/supervisor/conf.d/cron.conf
[program:cron]
command="cron"
process_name=%(program_name)s
autostart=true
autorestart=true
startretries=5
exitcodes=0,2,70
stopsignal=QUIT
stopwaitsecs=2
stdout_logfile=/var/log/cron.log
